# com.trisignal.trader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trisignal.trader</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/bytedance/.claude/skills/trisignal-trader/run_trisignal.sh</string>
    </array>

    <!-- 4H K线收盘后7分钟执行：北京时间 00:07/04:07/08:07/12:07/16:07/20:07 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>7</integer></dict>
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>7</integer></dict>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>7</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>7</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>7</integer></dict>
        <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>7</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/bytedance/.claude/skills/trisignal-trader/logs/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bytedance/.claude/skills/trisignal-trader/logs/launchd_err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/bytedance/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/bytedance</string>
    </dict>
</dict>
</plist>
